# crontab/cronjobs
#!/bin/bash
* * * * * /etc/cron.d/api-schedule.sh >> /var/log/cron.log 2>&1
